# examples/systemd/ksm-mcp.service
[Unit]
Description=KSM MCP Server - Keeper Secrets Manager Model Context Protocol
Documentation=https://github.com/keeper-security/ksm-mcp
After=network.target

[Service]
Type=simple
User=ksm-mcp
Group=ksm-mcp

# Working directory
WorkingDirectory=/var/lib/ksm-mcp

# Command
ExecStart=/usr/local/bin/ksm-mcp serve --batch --profile production

# Restart policy
Restart=always
RestartSec=10

# Environment
Environment="KSM_MCP_CONFIG_DIR=/etc/ksm-mcp"
Environment="KSM_MCP_LOG_DIR=/var/log/ksm-mcp"
Environment="KSM_MCP_PROFILE=production"
Environment="KSM_MCP_LOG_LEVEL=info"
Environment="KSM_MCP_BATCH_MODE=true"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ksm-mcp /var/log/ksm-mcp /etc/ksm-mcp

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
MemoryLimit=512M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ksm-mcp

[Install]
WantedBy=multi-user.target